# systemd/resourcemanager.socket
[Socket]
ListenStream=6312

[Install]
WantedBy=sockets.target
